# prometheus/rootfs/etc/services.d/prometheus-configgen/run
#!/command/with-contenv bashio
bashio::log.info 'Starting prometheus config generator...'

if ! bashio::fs.directory_exists /share/prometheus/targets; then
  mkdir -p /share/prometheus/targets
fi

cd /opt/prometheus-configgen || bashio::exit.nok \
  "Could not change working directory for prometheus-configgen."

# Run Prometheus
exec python3 combiner
